# worker/supervisord.conf
[supervisord]
logfile=/app/log/supervisord/supervisord.log    ; supervisord log file
logfile_maxbytes=50MB                           ; maximum size of logfile before rotation
logfile_backups=10                              ; number of backed up logfiles
loglevel=error                                  ; info, debug, warn, trace
pidfile=/app/run/supervisord.pid                ; pidfile location
nodaemon=true                                   ; run supervisord as a daemon
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors
user=www-data                                   ; default user
childlogdir=/app/log/supervisord/               ; where child log files will live

[program:worker]
command=python /app/worker.py
process_name=%(program_name)s-%(process_num)s
numprocs=5
directory=/app
stopsignal=TERM
autostart=true
autorestart=true